# package/deb/DEBIAN/prerm
#!/bin/bash

systemctl stop ssh-ito-denwa || true
systemctl disable ssh-ito-denwa || true

#rm /etc/systemd/system/multi-user.target.wants/ssh-ito-denwa.service 

